# redhat/parsoid.ini
[group:parsoid]
programs  = parsoid

[program:parsoid]
command         = mw-nserve
user            = nobody
directory       = /var/tmp
redirect_stderr = true
stdout_logfile  = /var/log/parsoid.log