# redhat/parsoid.ini
[group:parsoid]
programs  = parsoid

[program:parsoid]
command         = node bin/server.js
user            = nobody
directory       = /usr/lib64/node_modules/parsoid
redirect_stderr = true
stdout_logfile  = /var/log/parsoid.log